# testing/softether/vpnbridge.initd
#!/sbin/openrc-run

description="Starts or Stops the Softether VPN Bridge"
procname="vpnbridge"

start(){
	ebegin "Starting Softether VPN Bridge"
	/usr/bin/vpnbridge start
	eend $?
}

stop(){
	ebegin "Stopping Softether VPN Bridge"
	/usr/bin/vpnbridge stop
	eend $?
}
